# deploy/systemd/seller-dashboard-backup.timer
[Unit]
Description=Refresh the Seller Dashboard daily restore point every 6 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true
Unit=seller-dashboard-backup.service

[Install]
WantedBy=timers.target
